# test/config/systemd/fx-caddy-compose.path
[Unit]
Description=Docker compose change watcher for fx-caddy

[Path]
PathChanged=/opt/fx-mediaserver/config/fx-caddy/docker-compose.yml
Unit=fx-caddy-compose-reload.service

[Install]
WantedBy=fx-caddy.service
RequiredBy=fx-caddy.service
