# test/config/systemd/fx-qbittorrent.service
[Unit]
Description=Fixture torrent client (fx-qbittorrent)
After=docker.service mediaserver-network.service wireguard.service
Requires=docker.service mediaserver-network.service
PartOf=mediaserver-network.service mediaserver.target
BindsTo=wireguard.service

[Service]
Type=simple
TimeoutStartSec=0
Restart=on-failure
RestartSec=10
ExecStartPre=/bin/sh -c 'for i in {1..30}; do docker inspect wireguard --format="{{.State.Health.Status}}" 2>/dev/null | grep -q healthy && exit 0; sleep 1; done; exit 1'
ExecStart=/usr/bin/docker compose -f /opt/fx-mediaserver/config/fx-qbittorrent/docker-compose.yml up --force-recreate fx-qbittorrent
ExecStop=/usr/bin/docker compose -f /opt/fx-mediaserver/config/fx-qbittorrent/docker-compose.yml stop fx-qbittorrent

[Install]
WantedBy=mediaserver.target
