# test/config/systemd/fx-qbittorrent.path
[Unit]
Description=Config watcher for fx-qbittorrent

[Path]
Unit=fx-qbittorrent.service

[Install]
WantedBy=fx-qbittorrent.service
RequiredBy=fx-qbittorrent.service
