# systemd/service-compose-reload.service.erb
[Unit]
Description=Restart docker compose service <%= service.name %>
After=<%= service.name %>.service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart <%= service.name %>.service
