# systemd/service-compose.path.erb
[Unit]
Description=Docker compose change watcher for <%= service.name %>

[Path]
PathChanged=<%= install_base %>/config/<%= service.name %>/docker-compose.yml
Unit=<%= service.name %>-compose-reload.service

[Install]
WantedBy=<%= service.name %>.service
RequiredBy=<%= service.name %>.service
